# conf/upstart.conf
#ubuntu@ip-172-31-9-197:~$ cat /etc/init/fig.conf
description "fig"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
chdir /home/ubuntu/sms-panda

script
  # Wait for docker to finish starting up first.
  FILE=/var/run/docker.sock
  while [ ! -e $FILE ] ; do
    inotifywait -t 2 -e create $(dirname $FILE)
  done
  /usr/local/bin/fig up
end script